# com.startworking.keydetector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.startworking.keydetector.plist — macOS LaunchAgent for StartWorking

  This file tells macOS to start the key detector automatically at every login.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.startworking.keydetector</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/pasquale_marzaioli/Desktop/Coding/start_working/key_detector.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardInPath</key>
    <string>/dev/null</string>

    <key>StandardOutPath</key>
    <string>/Users/pasquale_marzaioli/Library/Logs/StartWorking/key_detector.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/pasquale_marzaioli/Library/Logs/StartWorking/key_detector.log</string>
</dict>
</plist>
