# venere/com.venere.clapdetector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.venere.clapdetector.plist — macOS LaunchAgent for Venere

  This file tells macOS to start the clap detector automatically at every login.

  INSTALLATION (replace YOUR_USERNAME with your actual macOS username):

    1. Copy this file to ~/Library/LaunchAgents/:
         cp venere/com.venere.clapdetector.plist ~/Library/LaunchAgents/

    2. Edit the copy and replace every occurrence of YOUR_USERNAME with your username:
         nano ~/Library/LaunchAgents/com.venere.clapdetector.plist

    3. Register the agent with launchd (no reboot needed):
         launchctl load ~/Library/LaunchAgents/com.venere.clapdetector.plist

    4. To unload (stop auto-start):
         launchctl unload ~/Library/LaunchAgents/com.venere.clapdetector.plist

  LOGS:
    ~/Library/VenereLauncher/launchagent.log   ← LaunchAgent stdout/stderr
    ~/Desktop/AntiGravity/venere/clap_detector.log  ← clap detector output
-->
<plist version="1.0">
<dict>
    <!-- Unique identifier for this agent — do not change -->
    <key>Label</key>
    <string>com.venere.clapdetector</string>

    <!-- Command to run at login:
         Opens a temporary Terminal window, starts the clap detector in the
         background, then closes the window. The delay gives the audio
         subsystem time to initialize before Python tries to open the mic. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/osascript</string>
        <string>-e</string>
        <!-- [CUSTOMIZE] Replace YOUR_USERNAME in the two paths below -->
        <string>delay 15
tell application "Terminal"
    do script "nohup /usr/bin/python3 /Users/YOUR_USERNAME/Desktop/AntiGravity/venere/clap_detector.py > /Users/YOUR_USERNAME/Desktop/AntiGravity/venere/clap_detector.log 2>&1 &; sleep 2; exit"
    activate
end tell
delay 4
tell application "Terminal" to close front window</string>
    </array>

    <!-- Run this agent every time the user logs in -->
    <key>RunAtLoad</key>
    <true/>

    <!-- [CUSTOMIZE] Replace YOUR_USERNAME in the two paths below -->
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Desktop/AntiGravity/venere/launchagent.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Desktop/AntiGravity/venere/launchagent.log</string>
</dict>
</plist>
